# phantomjs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:phantomjs]
command=/phantomjs/bin/phantomjs --webdriver=4444